# supervisor-haproxy.conf
[program:haproxy]
command = run_haproxy 1 /etc/haproxy/sentinel.cfg
priority = 200
startsecs = 2
exitcodes = 0
