# supervisor-sentinel.conf
[program:sentinel]
command = run_sentinel /etc/sentinel/haproxy.yml /etc/sentinel/haproxy.tpl
priority = 100
startsecs = 2
exitcodes = 0
